# systemd/iot-updater.service
[Unit]
Description=Cockpit IoT Updater — sidecar HTTP server (127.0.0.1:8088)
Documentation=https://github.com/legopc/cockpit-iot-updater
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /var/lib/iot-updater/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Sandboxing — sidecar needs /var/lib/iot-updater, /var/tmp, /run, /tmp
# Note: /tmp must be writable — bootc status writes a temp file there for SELinux install_t test
# Note: CapabilityBoundingSet is intentionally omitted — bootc status needs to enter the host
#       mount namespace (/proc/1/ns/mnt) which requires capabilities beyond any explicit list.
#       The service runs as root; ProtectSystem+ProtectHome provide the actual confinement.
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/iot-updater /var/tmp /run /tmp

[Install]
WantedBy=multi-user.target
